# deploy/systemd/slotfeed-coordinator.service
[Unit]
Description=SLOTFEED Stream Coordinator
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/slotfeed/backend
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/opt/slotfeed/backend/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/slotfeed/backend/venv/bin/python -m app.services.stream_coordinator --redis-url redis://localhost:6379/0 --check-interval 30 --job-interval 5
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slotfeed-coordinator

# Resource limits
MemoryMax=300M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
